# com.user.cls-pipeline-enrich.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.cls-pipeline-enrich</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/jintianyouyu/projects/cls-pipeline/scripts/run_enrich.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jintianyouyu/projects/cls-pipeline</string>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>SessionCreate</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/jintianyouyu/projects/cls-pipeline/logs/enrich-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jintianyouyu/projects/cls-pipeline/logs/enrich-stderr.log</string>
</dict>
</plist>
